# config/chromekey.service
[Service]
Type=simple
ExecStart=/usr/local/bin/chromekey --use_led=NUML
